# stealth-browser/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render
autorestart=true
priority=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:fluxbox]
command=fluxbox
autorestart=true
priority=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
depends_on=xvfb

[program:x11vnc]
command=x11vnc -display :99 -forever -shared -rfbport 5900 -nopw -xkb
autorestart=true
priority=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
depends_on=xvfb

[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autorestart=true
priority=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
depends_on=x11vnc

[program:api]
command=python -m uvicorn stealth_browser.src.main:app --host 0.0.0.0 --port 8001 --log-level info
autorestart=true
priority=4
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
depends_on=xvfb
environment=DISPLAY=":99"
